# deploy/usocksvc/usocksvc.service
[Unit]
Description=Temperature Humidity socket service

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=thsock


